# build-output/kOS/includes.chroot/etc/skel/.config/autostart/set_wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Set Wallpaper
Exec=/usr/local/bin/kos-set-wallpaper
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
